# ami-proxy/ami-proxy.service
[Unit]
After=network.target

[Service]
Type=simple
ExecStart=/usr/sbin/amiproxy

[Install]
WantedBy=multi-user.target
